# systemd/tunet-bash.timer
[Unit]
Description=Tsinghua University Network Authentication Timer

[Timer]
OnStartupSec=1m
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target
